# deploy/systemd/argus-auto-deploy.service
[Unit]
Description=Argus auto deploy from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=argus
Group=argus
WorkingDirectory=/opt/argus
TimeoutStartSec=30min
ExecStart=/usr/local/bin/argus-run-background-job.sh argus-auto-deploy /usr/local/bin/argus-auto-deploy.sh
