# deploy/systemd/argus-unread-reminders.timer
[Unit]
Description=Run Argus unread reminders every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
Unit=argus-unread-reminders.service

[Install]
WantedBy=timers.target
