# deploy/systemd/argus-cleanup-old-leads.service
[Unit]
Description=Argus cleanup old leads

[Service]
Type=oneshot
User=argus
Group=argus
WorkingDirectory=/opt/argus
Environment="PATH=/opt/argus/.venv/bin"
TimeoutStartSec=30min
ExecStart=/usr/local/bin/argus-run-background-job.sh argus-cleanup-old-leads /opt/argus/.venv/bin/python manage.py cleanup_old_leads --days 30 --limit 100
